# goku-lite.service
[Unit]
Description=Goku Lite AI Orchestrator
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/afrifa-gilbert/Documents/Dev/goku/goku-lite
Environment="PATH=/home/afrifa-gilbert/Documents/Dev/goku/goku-lite/venv/bin"
ExecStart=/home/afrifa-gilbert/Documents/Dev/goku/goku-lite/venv/bin/python main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
